# scripts/tradingbot.service
# PolyPrinting2026 - Systemd Service File
#
# Installation:
#   sudo cp tradingbot.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable tradingbot
#   sudo systemctl start tradingbot
#
# Management:
#   sudo systemctl status tradingbot
#   sudo systemctl stop tradingbot
#   sudo systemctl restart tradingbot
#   sudo journalctl -u tradingbot -f

[Unit]
Description=PolyPrinting2026 - Kalshi/Polymarket Trading Bot
Documentation=https://github.com/YOUR_USERNAME/PolyPrinting2026
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu

# Working directory (update path as needed)
WorkingDirectory=/home/ubuntu/trading-bot

# Environment file with API keys
# For Kalshi (US): use kalshi_env
# For Polymarket (non-US): use polymarket_env
EnvironmentFile=/home/ubuntu/.kalshi_env

# Python virtual environment activation and bot execution
# Change --exchange to 'polymarket' for non-US users
ExecStart=/home/ubuntu/trading-bot/venv/bin/python main.py --exchange kalshi --live

# Restart policy
Restart=always
RestartSec=10

# Stop signal (allow graceful shutdown)
KillSignal=SIGTERM
TimeoutStopSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tradingbot

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes

# Allow write access to logs directory
ReadWritePaths=/home/ubuntu/trading-bot/logs
ReadWritePaths=/home/ubuntu/trading-bot/data

# Resource limits (adjust based on your instance)
MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
